# ops/systemd/iic-promoter.service
# /etc/systemd/system/iic-promoter.service
#
# IIC-FORGE F4 — promoter service
# Polls events table every 10s, enqueues event_alert jobs into queue_jobs.
# Defensive retry-internal: see tradingagents/orchestrator/promoter.py
# The exit-gate evaluator (scripts/f4_exit_gate.py) fails the gate if this
# unit's NRestarts > 0 during the observation window.

[Unit]
Description=IIC-FORGE — F4 promoter (events → queue_jobs)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ziwei-huang
WorkingDirectory=/home/ziwei-huang/TradingAgents/TradingAgents
EnvironmentFile=/home/ziwei-huang/TradingAgents/TradingAgents/.env
# ── Local LLM cutover (IIC-FORGE-05) ─────────────────────────────────────────
# To route the alert_gate role to a local llama-server instance, add the four
# lines below to the EnvironmentFile (.env) and restart the unit
# (daemon-reload is only needed when unit files change, not for .env edits):
#   sudo systemctl restart iic-promoter
#
# Cutover (uncomment + fill in .env, or add as Environment= lines here):
#   LOCAL_LLM_BASE_URL=http://127.0.0.1:8080/v1
#   LOCAL_LLM_API_KEY=<bearer-token if llama-server was started with --api-key; omit otherwise>
#   IIC_ALERT_GATE_LLM_PROVIDER=local
#   IIC_ALERT_GATE_LLM_MODEL=<model-id matching what llama-server reports, e.g. qwen3.6-27b-instruct-q4_k_m>
#
# Revert: remove (or comment out) ALL FOUR IIC_ALERT_GATE_* vars (both provider
# and model). Leaving IIC_ALERT_GATE_LLM_MODEL set while unsetting the provider
# yields provider=deepseek + model=qwen…, causing API model-not-found errors.
# LOCAL_LLM_BASE_URL and LOCAL_LLM_API_KEY are safe to leave set; they are only
# consumed when provider=local is active.
# ─────────────────────────────────────────────────────────────────────────────
ExecStart=/home/ziwei-huang/miniconda3/bin/python -m tradingagents.orchestrator.promoter
Restart=on-failure
RestartSec=30
MemoryMax=256M
CPUQuota=25%
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
