# ops/systemd/iic-triage.service
[Unit]
Description=IIC triage consumer — dedupe + salience + persistence
After=network-online.target redis-server.service
Requires=redis-server.service

[Service]
Type=simple
User=ziwei-huang
WorkingDirectory=/home/ziwei-huang/TradingAgents/TradingAgents
EnvironmentFile=/home/ziwei-huang/TradingAgents/TradingAgents/.env
# ── Local LLM cutover (IIC-FORGE-05) ─────────────────────────────────────────
# To route the triage_salience role to a local llama-server instance, add the
# four lines below to the EnvironmentFile (.env) and restart the unit
# (daemon-reload is only needed when unit files change, not for .env edits):
#   sudo systemctl restart iic-triage
#
# Cutover (uncomment + fill in .env, or add as Environment= lines here):
#   LOCAL_LLM_BASE_URL=http://127.0.0.1:8080/v1
#   LOCAL_LLM_API_KEY=<bearer-token if llama-server was started with --api-key; omit otherwise>
#   IIC_TRIAGE_LLM_PROVIDER=local
#   IIC_TRIAGE_LLM_MODEL=<model-id matching what llama-server reports, e.g. qwen3.6-27b-instruct-q4_k_m>
#
# Revert: remove (or comment out) ALL FOUR IIC_TRIAGE_* vars (both provider
# and model). Leaving IIC_TRIAGE_LLM_MODEL set while unsetting the provider
# yields provider=deepseek + model=qwen…, causing API model-not-found errors.
# LOCAL_LLM_BASE_URL and LOCAL_LLM_API_KEY are safe to leave set; they are only
# consumed when provider=local is active.
# ─────────────────────────────────────────────────────────────────────────────
ExecStart=/home/ziwei-huang/miniconda3/bin/python -m tradingagents.sensing.triage
Restart=on-failure
RestartSec=30
MemoryMax=2G
CPUQuota=100%
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
