# ops/systemd/redis-server.service
[Unit]
Description=Legacy placeholder - Redis is owned by IIC-Forge Compose

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo "Redis is owned by compose.yml service redis; do not start a standalone legacy container."'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
